# sensor.service
[Unit]
Description=Sensor detecting and distribution service

[Service]
WorkingDirectory=/home/pi/bin/DepotKontroll/
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/pi/bin/DepotKontroll/.env
ExecStart=/usr/bin/python3 /home/pi/bin/DepotKontroll/readings.py

Restart=on-failure

[Install]
WantedBy=default.target
